# buildpacks/print-message/layer/print-message/profile.d/message.sh
export MESSAGE="$(cat Message)"
export OTHER_MESSAGE="Yes, we found a Message file. See \$MESSAGE for contents."